# jenkins.ini
[program:jenkins]
command=/usr/bin/java -jar /jenkins/jenkins.war
user=jenkins
autostart=true
autorestart=true
environment=HOME=/var/lib/jenkins
priority=2
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
